# scripts/supervisord.conf
[supervisord]
nodaemon=true

[program:rabbit]
command=/bin/bash -c "docker-entrypoint.sh rabbitmq-server"

[program:mongo]
command=/bin/bash -c "/usr/local/bin/docker-mongo-entrypoint.sh mongod"